# paketbote-check-aur-versions.timer
[Unit]
Description = Paketbote Check AUR Versions

[Timer]
OnCalendar = daily

[Install]
WantedBy = multi-user.target

